# openhomebus/rootfs/etc/services.d/ohb-dashboard/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# OpenHomeBus Dashboard — S6 finish script (clean shutdown)
# ==============================================================================

bashio::log.info "OpenHomeBus Dashboard stopped"
